# Mixed/LogitechMouseFix/g502.rules
ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c07d", RUN+="/run/media/sdar/Storage_2/config/scripts/g502/reconnect.sh"
